# files/system/shared/usr/lib/systemd/system/birralee-update.service
[Unit]
Description=Birralee auto update

[Service]
Type=oneshot
ExecStart=/usr/libexec/birralee/birralee-update
