# files/system/shared/usr/lib/systemd/system/flatpak-update.service
[Unit]
Description=Flatpak auto update

[Service]
Type=oneshot
ExecStart=/usr/libexec/birralee/flatpak-update
